# config/conf.d/program2.conf
[program:dummy2]
command=/bin/sleep 1000
autostart=true
autorestart=true
stderr_logfile=/home/kobay/dummy.err.log
stdout_logfile=/home/kobay/dummy.out.log
