# Taskmasterd/etc/conf.d/program1.conf

[program:hostname]
command=hostname
user=1000
