# Taskmaster/program2.ini
[program:dummy2]
command=/bin/sleep 1000
autostart=true
autorestart=trues
stderr_logfile=/home/vzurera-/dummy.err.log
stdout_logfile=/home/vzurera-/dummy.out.log
